# kernel_logs.timer
[Unit]
Description=Timer logs hourly

[Timer]
OnUnitActiveSec=1h
Persistent=true

[Install]
WantedBy=timers.target
